# dbus/machine-learning-agent.service.in
[Unit]
Description=Machine Learning Agent Daemon

[Service]
Type=simple
SmackProcessLabel=System
ExecStart=/usr/bin/machine-learning-agent
