# dbus/machine-learning-agent.service.in
[Unit]
Description=Machine Learning Agent Daemon

[Service]
Type=dbus
BusName=org.tizen.machinelearning.service
SmackProcessLabel=System
ExecStart=/usr/bin/machine-learning-agent
User=service_fw
Group=service_fw
